# ansible/templates/deploy.macmini.service.j2
# FastDeploy service for macmini
# /etc/systemd/system/fastdeploy.service
[Unit]
Description=FastDeploy Service on macmini
After=network.target postgresql.service
Requires=postgresql.service

[Service]
Type=simple
Restart=always
RestartSec=10
WorkingDirectory={{ site_path }}
User={{ username }}
Group={{ username }}

# Environment variables
Environment="PATH={{ venv_bin }}:/usr/local/bin:/usr/bin:/bin"
EnvironmentFile={{ site_path }}/.env

# Start command
ExecStart={{ site_path }}/.venv/bin/python -m uvicorn deploy.entrypoints.fastapi_app:app --port {{ app_port }} --host 127.0.0.1

# Process management
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Progressive security hardening - Phase 1
# These settings provide basic security without breaking functionality
PrivateTmp=yes
NoNewPrivileges=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes

# Filesystem protections - moderate level
ProtectSystem=full
ProtectHome=read-only

# Network - allow all needed families
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX

# Memory limits
MemoryMax=2G

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=fastdeploy

[Install]
WantedBy=multi-user.target
